# tools-for-vps/grab-ntp-time.service
[Unit]
Description=Update time via NTP

[Service]
Type=oneshot
User=root
ExecStart=/root/tools/bin/ntpd -ddgqc /dev/null fr.pool.ntp.org
