# tools-for-vps/ii-oftc.service
[Unit]
Description=ii IRC watcher

[Service]
User=raskin
ExecStart=/home/raskin/ii-starter-oftc

[Install]
WantedBy=multi-user.target
