# tools-for-vps/ii.service
[Unit]
Description=ii IRC watcher

[Service]
User=raskin
ExecStart=/home/raskin/ii-starter

[Install]
WantedBy=multi-user.target
