# tools-for-vps/openvpn-tcp.service
[Unit]
Description=OpenVPN server (UDP)

[Service]
ExecStart=/root/tools/bin/openvpn /etc/openvpn/server-tcp.conf
WorkingDirectory=/etc/openvpn/

[Install]
WantedBy=multi-user.target
