# tools-for-vps/postfix-key-concat.timer
[Unit]
Description=Run key concatenation daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent = True

[Install]
WantedBy=multi-user.target

